# infoburp.conf
#!upstart
description "infoburp"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

exec node ~/infoburp/server.js >> /var/log/infoburp.log 2>&1
